# systemd/usr__lib__systemd__system__bootc-fetch.service
[Unit]
Description=Fetch bootc image
After=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/bin/bootc update
