# systemd/usr__lib__systemd__system__bootc-fetch.timer
[Unit]
Description=Fetch bootc image daily
ConditionPathExists=/run/ostree-booted

[Timer]
OnCalendar=*-*-* 12:30:00
Persistent=true

[Install]
WantedBy=timers.target
