# systemd/usr__lib__systemd__system__bootc-fetch.timer
[Unit]
Description=Fetch bootc image daily

[Timer]
OnCalendar=*-*-* 12:30:00
Persistent=true

[Install]
WantedBy=timers.target
